# supervisor/example-dotnet.conf
[program:example-dotnet]
directory=/opt/example-dotnet
command=sh -c 'dotnet restore && dotnet run'
